# templates/app-service/iota-inx-poi.service
# Version 0
[Unit]
Description=IOTA-INX-POI
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/default/iota-inx-poi
LimitNOFILE=8192
SyslogIdentifier=iota-inx-poi
PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes

User=inx
WorkingDirectory=/var/lib/iota-inx-poi
TimeoutSec=1200
Restart=always
ExecStart=/usr/bin/iota-inx-poi $OPTIONS

[Install]
WantedBy=multi-user.target